# roles/ripper/templates/brain.service.j2
[Service]
ExecStart={{ binroot }}/brain.py

