# roles/ripper/templates/brain.service.j2
[Service]
User={{ ansible_user }}
ExecStart={{ binroot }}/brain.py

[Install]
WantedBy=multi-user.target
